# container/service/php-fpm7/run
#!/bin/execlineb
fdmove -c 2 1  # stderr to stdout so it is logged neatly
php-fpm7 -F -O
