# container/service/php-fpm7/run
#!/bin/execlineb
fdmove -c 2 1 php-fpm7 -F -O
