# container/service/schema/finish
#!/bin/execlineb
exit 125  # run once only
