# container/service/nginx/run
#!/bin/execlineb
nginx -c /opt/nginx.conf
