# container/service/nginx/run
#!/bin/execlineb
fdmove -c 2 1 nginx -c /opt/nginx.conf
