# container/service/nginx/run
#!/bin/execlineb
fdmove -c 2 1  # stderr to stdout so it is logged neatly
nginx -c /opt/nginx.conf
