# container/service/mariadb/run
#!/bin/execlineb
fdmove -c 2 1 cd /opt/data /usr/bin/mysqld_safe --datadir=/opt/data
