# container/service/mariadb/run
#!/bin/execlineb
cd /opt/data /usr/bin/mysqld_safe --datadir=/opt/data
