# container/service/db-setup/finish
#!/bin/execlineb
exit 125  # We only intended to run once; don't restart us
